# KaguraMeaLive-subscribe.timer
[Unit]
Description=subscribe daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
